# systemd/english-companion-nx.service
[Unit]
Description=English Companion NX - AI Conversational English Practice
Documentation=https://github.com/robertkotsch/english-companion-nx
After=network.target ollama.service

[Service]
Type=simple
WorkingDirectory=%h/apps/english-companion-nx
Environment="PATH=%h/apps/english-companion-nx/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="VIRTUAL_ENV=%h/apps/english-companion-nx/.venv"

# Run the conversation prototype
ExecStart=%h/apps/english-companion-nx/.venv/bin/python conversation_prototype.py

# Resource Limits (Jetson Orin NX 16GB - allocate 11GB max)
MemoryMax=11G
MemoryHigh=10G

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=300

# Standard output/error to journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=english-companion-nx

# Nice level (slightly higher priority)
Nice=-5

# Security (optional hardening)
# PrivateTmp=yes
# NoNewPrivileges=yes

[Install]
WantedBy=default.target
